# pysel_scoring.service
[Unit]
Description=PySEL Scoring Service

[Service]
Type=simple
ExecStartPre=-/bin/sleep 45
ExecStart=/usr/bin/python3 /usr/local/bin/dist/score.py

[Install]
WantedBy=multi-user.target
