# static/pysel_scoring.service
[Unit]
Description=PySEL Scoring Service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/bin/PySEL.py

[Install]
WantedBy=multi-user.target
